# overlay/etc/systemd/system/nginx.service.d/restart.conf
[Service]
# Auto-restart settings
Restart=on-failure
RestartSec=5